# scripts/debian/postinst
if ! [ -f "/opt/report-engine/keys/report-engine.key" ]; then
    mkdir -p '/opt/report-engine/keys/'
    openssl req -x509 -newkey rsa:4096 -subj "/CN=$(hostname -I | cut -d" " -f1 | xargs)" -addext "subjectAltName=IP:$(hostname -I | cut -d" " -f1 | xargs),IP:127.0.0.1,DNS:$(hostname)" -keyout /opt/report-engine/keys/report-engine.key -nodes -out /opt/report-engine/keys/report-engine.pem -sha256 -days 358000
fi
        
chmod -R 770 /opt/report-engine
chown -R root:root /opt/report-engine

if [ -f "/usr/lib/systemd/system/report-engine.service" ]; then
    rm -rf /usr/lib/systemd/system/report-engine.service  2>/dev/null || true
    systemctl disable report-engine.service 2>/dev/null || true
    systemctl stop report-engine.service 2>/dev/null || true
    systemctl daemon-reload 2>/dev/null || true
fi

echo """
[Unit]
Description=FastAPI Report Engine
[Service]
Type=simple
WorkingDirectory=/opt/report-engine
ExecStart=/opt/report-engine/report-engine
Restart=always
RestartSec=10
SyslogIdentifier=report-engine
KillSignal=SIGINT
User=root
Group=root
[Install]
WantedBy=multi-user.target
    """ > /etc/systemd/system/report-engine.service

systemctl daemon-reload
systemctl enable report-engine.service
systemctl restart report-engine.service